# packages/techmago-settings/rhel8/workdir/etc/skel/.config/autostart/gdm-skell.desktop
[Desktop Entry]
Name=Techmago initial gdm config
Icon=preferences-system
Exec=/usr/libexec/techload-gdm
Terminal=true
Type=Application
StartupNotify=true
Categories=GNOME;GTK;System;
OnlyShowIn=GNOME;
NoDisplay=false
